# scripts/supervisord.conf
[supervisord]
nodaemon=true
logfile=/data/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
pidfile=/var/run/supervisord.pid

[program:access_server]
command=python -m backend.services.access_server
directory=/app
autorestart=true
priority=5
startsecs=3
stdout_logfile=/data/logs/access_server.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=2
stderr_logfile=/data/logs/access_server_err.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=2

[program:postfix]
command=/app/scripts/postfix-start.sh
autorestart=true
priority=10
stdout_logfile=/data/logs/postfix.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/data/logs/postfix_err.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=2

[program:fastapi]
command=uvicorn backend.main:app --host 0.0.0.0 --port 8080
directory=/app
autorestart=true
priority=20
stdout_logfile=/data/logs/fastapi.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/data/logs/fastapi_err.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=2

[program:health_checker]
command=python -m backend.services.health_checker
directory=/app
autorestart=true
priority=30
startsecs=10
stdout_logfile=/data/logs/health.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=2
stderr_logfile=/data/logs/health_err.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=2

[program:mail_logger]
command=python -m backend.services.mail_logger
directory=/app
autorestart=true
priority=30
stdout_logfile=/data/logs/mail_logger.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=2
stderr_logfile=/data/logs/mail_logger_err.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=2

[program:expiry_checker]
command=python -m backend.services.expiry_checker
directory=/app
autorestart=true
priority=30
startsecs=10
stdout_logfile=/data/logs/expiry.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=2
stderr_logfile=/data/logs/expiry_err.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=2

[program:log_cleanup]
command=python -m backend.services.log_cleanup
directory=/app
autorestart=true
priority=40
startsecs=5
stdout_logfile=/data/logs/log_cleanup.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=1
stderr_logfile=/data/logs/log_cleanup_err.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=1
